# 10_updater/overlay/usr/lib/systemd/system/neon-updater.service
[Unit]
Description=Neon Updater

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/neon/versions.conf
ExecStart=/opt/neon/update $NEON_CORE_REF
