# 05_neon_core/overlay/usr/lib/systemd/system/neon-enclosure.service
[Unit]
Description=Neon Enclosure
PartOf=neon.service
After=neon.service
After=neon-bus.service
After=neon-gui.service

[Service]
Type=notify
TimeoutStopSec=30
Restart=always
User=neon
EnvironmentFile=/etc/neon/neon_env.conf
WorkingDirectory=/opt/neon
ExecStart=/usr/libexec/neon-enclosure.py
#StandardOutput=append:/home/neon/logs/enclosure.log
#StandardError=inherit
[Install]
WantedBy=neon.service
