# 05_neon_core/overlay/usr/lib/systemd/system/neon-bus.service
[Unit]
Description=Neon Message Bus
PartOf=neon.service
After=neon.service
After=neon-logs.service

[Service]
Type=notify
Restart=always
User=neon
WorkingDirectory=/opt/neon
# Start module when unit is started
ExecStart=/usr/libexec/neon-messagebus.py
# Stop module when unit is stopped
#ExecStop=/bin/bash /opt/neon/kill_service.sh neon_messagebus_service
StandardOutput=append:/home/neon/logs/bus.log
StandardError=inherit
[Install]
WantedBy=neon.service
