# base_network_manager/overlay/usr/lib/systemd/system/wifi-setup.service
[Unit]
Description=Wifi Setup watcher
Requires=network-manager

[Service]
ExecStart=/bin/bash /opt/neon/check_wifi_setup.sh

[Install]
WantedBy=multi-user.target
