# 11_factory_reset/overlay/usr/lib/systemd/system/neon-reset.service
[Unit]
Description=Neon Reset

[Service]
Type=oneshot
User=root
ExecStart=/opt/neon/reset
